# systemd/central-key-manager.service.in
[Unit]
Description=Start the Central Key Manager
DefaultDependencies=no

[Service]
User=key-manager
Group=key-manager
SmackProcessLabel=System
Type=notify
ExecStart=/usr/bin/key-manager
Sockets=central-key-manager-api-storage.socket
Sockets=central-key-manager-api-control.socket
Sockets=central-key-manager-api-ocsp.socket
Sockets=central-key-manager-api-encryption.socket
EnvironmentFile=-@SYSTEMD_ENV_FILE@
RuntimeDirectory=@SERVICE_NAME@

[Install]
WantedBy=multi-user.target
